# trisul-probe/oink.cron
0 0 * * * oinkmaster -C /etc/oinkmaster.conf -o /etc/suricata/rules
